# airootfs/etc/s6/sv/pacman-init/up
#!/usr/bin/execlineb -P

fdmove -c 2 1
foreground { pacman-key --init }
foreground { pacman-key --populate artix }
pacman-key --populate archlinux
